# packer/bind-dns/assets/systemd/named.service.d/override.conf
[Unit]
Description=BIND recursive resolver for aws-firewall-proxy shared DNS cache

[Service]
Restart=on-failure
RestartSec=5s
